# backend/worker/supervisor/simhpc.conf
[supervisord]
nodaemon=true
user=appuser

[program:api]
command=python -m backend
directory=/app/backend
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
priority=10

[program:worker]
command=python -m backend.worker.worker
directory=/app/backend
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
priority=20

[program:autoscaler]
command=python -m backend.runtime.autoscaler  ; adjust if you create this
directory=/app/backend
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
priority=30
