# docker/supervisor/simhpc.conf
[supervisord]
nodaemon=true

[program:api]
command=uvicorn app.api.main:app --host 0.0.0.0 --port 8080
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:worker]
command=python -m worker.runtime.runner
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:autoscaler]
command=python -m autoscaler.main
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:auto_stop_reaper]
command=python backend/scripts/auto_stop.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr