# docker/supervisor/simhpc.conf
[supervisord]
nodaemon=true

[program:api]
# Points to gateway.py inside backend/
command=python3 -m uvicorn backend.app.gateway:app --host 0.0.0.0 --port 8080
directory=/app
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:worker]
# Matches the package structure in your project
command=python3 -m backend.worker.worker
directory=/app
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:autoscaler]
command=python3 -m backend.autoscaler.main
directory=/app
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
