# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
user=root

[program:api]
command=python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8080
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
# Wait for API and Redis before starting worker
command=sh -c "python3 docker/scripts/wait-for-redis.py && sh docker/scripts/wait-for-api.sh && python3 app/services/worker/worker.py"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:autoscaler]
# Wait for API and Redis before starting autoscaler
command=sh -c "python3 docker/scripts/wait-for-redis.py && sh docker/scripts/wait-for-api.sh && python3 app/services/worker/autoscaler.py"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
